# docs/deploy/bot.service
[Unit]
Description=Omni Bot (Telegram long polling)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/omni/apps/bot
EnvironmentFile=/etc/omni/bot.env
ExecStart=/root/.bun/bin/bun run dev
Restart=always
RestartSec=3
TimeoutStartSec=30
TimeoutStopSec=30
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
